# etc/profile.d/motd.sh
#!/usr/bin/env bash

echo "TArmbian 26.5.1 (trixie)"
echo "Made by bangkkuser"
echo "Be welcome!"
echo
